# system-config/.config/systemd/user/fs-index.timer
[Unit]
Description=Timer to index filesystem for fzf

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
