# system-config/.config/systemd/user/fs-index.service
[Unit]
Description=Index filesystem for fzf

[Service]
Type=simple
ExecStart=/home/cdchawthorne/bin/fs
