# system-config/.config/systemd/user/get-mail.service
[Unit]
Description=Fetch mail via IMAP

[Service]
Type=oneshot
ExecStart=mbsync -a
ExecStart=notmuch new
